# deploy/freqtrade-monitor.service
[Unit]
Description=Freqtrade monitor dashboard
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/freqtrade-strategies
EnvironmentFile=/home/ubuntu/freqtrade-strategies/user_data/monitor.env
ExecStart=/usr/bin/node /home/ubuntu/freqtrade-strategies/dashboard/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
