# deploy/freqtrade-data-reliability.timer
[Unit]
Description=Run Freqtrade data reliability checks every five minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=freqtrade-data-reliability.service

[Install]
WantedBy=timers.target
