# research/runtime/systemd/research-registry-backup.timer
[Unit]
Description=Daily research governance registry verified backup

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=15m
Persistent=true
Unit=research-registry-backup.service

[Install]
WantedBy=timers.target
